# bot.conf
[program:alleycatbot]
directory=/home/ubuntu/AlleyCat-bot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME='/home/ubuntu/',USER='ubuntu'

stderr_logfile=/home/ubuntu/AlleyCat-bot/logfile_err.log
stdout_logfile=/home/ubuntu/AlleyCat-bot/logfile.log
